# systemd/pidrive_web.service
[Unit]
Description=PiDrive Web UI (Debug/Remote)
After=network-online.target
Wants=network-online.target
# KEIN After=pidrive_core — erzeugt Ordering-Cycle beim Shutdown
# Web läuft unabhängig und zeigt "Core offline" wenn nötig
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pidrive/pidrive
ExecStart=/usr/bin/python3 /home/pi/pidrive/pidrive/webui.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
