# systemd/pidrive_display.service
[Unit]
Description=PiDrive Display - Statusanzeige
# GPT-5.4: kein Requires — Display-Fehler stoppen nicht den Core
After=multi-user.target
# KEIN After=pidrive_core — vermeidet Ordering-Cycle
# Nur starten wenn SPI-Display vorhanden (fb1)
ConditionPathExists=/sys/class/graphics/fb1
# max 3 Neustarts in 120s
StartLimitIntervalSec=300
StartLimitBurst=2
RestartSec=30

[Service]
UMask=0000
Type=simple
User=root

# Direkt auf SPI-Display fb1 — kein fbcp noetig
Environment=SDL_FBDEV=/dev/fb1
Environment=SDL_VIDEODRIVER=fbcon
Environment=SDL_NOMOUSE=1
Environment=SDL_AUDIODRIVER=dummy
Environment=SDL_VIDEO_FBCON_KEEP_TTY=1
WorkingDirectory=/home/pi/pidrive/pidrive

ExecStartPre=/bin/sleep 5
ExecStartPre=/bin/bash -c 'cat /sys/class/graphics/fb1/virtual_size 2>/dev/null | grep -q "^[0-9]" || exit 77'
ExecStartPre=/bin/bash -c 'echo 0 > /sys/class/vtconsole/vtcon1/bind || true'
ExecStartPre=/bin/bash -c 'echo 0 > /sys/class/graphics/fbcon/cursor_blink || true'
ExecStart=/usr/bin/python3 /home/pi/pidrive/pidrive/main_display.py

# GPT-5.4: Crash neu starten, aber mit Limit
Restart=on-failure
RestartSec=10
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
