# systemd/pidrive.service
[Unit]
Description=PiDrive - Car Infotainment
After=multi-user.target rc-local.service
Conflicts=getty@tty1.service getty@tty2.service getty@tty3.service

[Service]
Type=simple
User=root

# Kein PAMName/TTYPath/StandardInput - haengt systemd 247 PAM-Helper
# setsid+TIOCSCTTY in launcher.py - funktioniert seit v0.5.3

Environment=SDL_FBDEV=/dev/fb0
Environment=SDL_VIDEODRIVER=fbcon
Environment=SDL_NOMOUSE=1
Environment=SDL_AUDIODRIVER=dummy
WorkingDirectory=/home/pi/pidrive/pidrive

ExecStartPre=/bin/sleep 3
ExecStart=/usr/bin/python3 /home/pi/pidrive/pidrive/launcher.py

Restart=always
RestartSec=10
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
