# systemd/pidrive_core.service
[Unit]
Description=PiDrive Core - Car Infotainment
After=multi-user.target network-online.target pulseaudio.service bluetooth.service
Wants=network-online.target pulseaudio.service bluetooth.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
UMask=0000
Type=simple
User=root
WorkingDirectory=/home/pi/pidrive/pidrive

# v0.9.0: PULSE_SERVER sicherstellen - alle Kindprozesse (mpv, rtl_fm|mpv, welle-cli|mpv)
# nutzen damit den System-PulseAudio-Socket statt User-Session
Environment=PULSE_SERVER=unix:/var/run/pulse/native

ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/python3 /home/pi/pidrive/pidrive/main_core.py

Restart=on-failure
RestartSec=5
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
