# systemd/pidrive_avrcp.service
[Unit]
Description=PiDrive AVRCP Trigger (BMW iDrive → File-Trigger)
After=bluetooth.target
Wants=bluetooth.target
# KEIN After=pidrive_core.service — vermeidet systemd Ordering-Cycle beim Shutdown (v0.8.7)

[Service]
UMask=0000
Type=simple
User=root
WorkingDirectory=/home/pi/pidrive/pidrive
ExecStart=/usr/bin/python3 /home/pi/pidrive/pidrive/avrcp_trigger.py
Restart=always
RestartSec=5
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
