# systemd/pipewire.service
[Unit]
Description=PiDrive Audio Server (PipeWire, System Mode)
Documentation=https://pipewire.org
After=dbus.service
Before=pipewire-pulse.service wireplumber.service

[Service]
Type=simple
User=pulse
Group=pulse
SupplementaryGroups=audio bluetooth
RuntimeDirectory=pipewire
RuntimeDirectoryMode=0755
Environment=XDG_RUNTIME_DIR=/run/pipewire
Environment=PIPEWIRE_RUNTIME_DIR=/run/pipewire
# Session-Bus auf System-Bus zeigen — PipeWire braucht spa.dbus für BT/ALSA
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket
ExecStart=/usr/bin/pipewire
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
