# systemd/wireplumber.service
[Unit]
Description=PiDrive Session Manager (WirePlumber, System Mode)
After=pipewire.service bluetooth.service
Requires=pipewire.service
Wants=bluetooth.service

[Service]
Type=simple
User=pulse
Group=pulse
SupplementaryGroups=audio bluetooth
Environment=XDG_RUNTIME_DIR=/run/pipewire
Environment=PIPEWIRE_RUNTIME_DIR=/run/pipewire
# Session-Bus auf System-Bus zeigen — WirePlumber braucht D-Bus für BT-Plugin
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket
Environment=HOME=/var/lib/pipewire
ExecStart=/usr/bin/wireplumber -p main-systemwide
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
